# API/supervisor/externalentropy.conf
[program:externalentropy]
autorestart = True
directory = /home/pi/
user = pi
command = /home/pi/.distributedurandom/externalentropy.py
environment = STNORESTART="1", HOME="/home/pi"
